# src/sysv/systemd/sssd.service.in
[Unit]
Description=System Security Services Daemon
# If dbus crashes or is stopped, then SSSD will be stopped, too
Requires=dbus.socket
# SSSD will not be started until syslog is
After=syslog.target

[Service]
ExecStart=@sbindir@/sssd -D -f
# These two should be used with traditional UNIX forking daemons
# consult systemd.service(5) for more details
Type=forking
PIDFile=@localstatedir@/run/sssd.pid

[Install]
WantedBy=multi-user.target
