# src/sysv/systemd/sssd-ssh.socket.in
[Unit]
Description=SSSD SSH Service responder socket
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
DefaultDependencies=no
Conflicts=shutdown.target

[Socket]
ExecStartPre=@libexecdir@/sssd/sssd_check_socket_activated_responders -r ssh
ListenStream=@pipepath@/ssh
SocketUser=@SSSD_USER@
SocketGroup=@SSSD_USER@

[Install]
WantedBy=sssd.service
