# src/sysv/systemd/sssd-kcm.service.in
[Unit]
Description=SSSD Kerberos Cache Manager
Documentation=man:sssd-kcm(5)
Requires=sssd-kcm.socket
After=sssd-kcm.socket

[Install]
Also=sssd-kcm.socket

[Service]
Environment=DEBUG_LOGGER=--logger=files
ExecStartPre=+-/bin/chown -f root:root @sssdconfdir@
ExecStartPre=+-/bin/chown -f root:root @sssdconfdir@/sssd.conf
ExecStartPre=+-/bin/chown -f -R root:root @sssdconfdir@/conf.d
ExecStartPre=+-find -type f -iname @secdbpath@/*.ldb -exec /bin/chown -f root:root {};
ExecStartPre=+-/bin/chown -f root:root @logpath@/sssd_kcm.log
ExecStart=@libexecdir@/sssd/sssd_kcm ${DEBUG_LOGGER}
CapabilityBoundingSet= CAP_DAC_OVERRIDE CAP_CHOWN CAP_SETGID CAP_SETUID
SecureBits=noroot noroot-locked
User=root
Group=root
@supplementary_groups@
