# src/sysv/systemd/sssd-kcm.socket.in
[Unit]
Description=SSSD Kerberos Cache Manager responder socket
Documentation=man:sssd-kcm(8)
Requires=sssd-secrets.socket

[Socket]
ListenStream=@localstatedir@/run/.heim_org.h5l.kcm-socket

[Install]
WantedBy=sockets.target
