# src/sysv/systemd/sssd-sudo.socket.in
[Unit]
Description=SSSD Sudo Service responder socket
Documentation=man:sssd.conf(5)

[Socket]
ListenStream=@pipepath@/sudo

[Install]
WantedBy=sockets.target
