# src/sysv/systemd/sssd.service.in
[Unit]
Description=System Security Services Daemon
# SSSD must be running before we permit user sessions
Before=systemd-user-sessions.service nss-user-lookup.target
Wants=nss-user-lookup.target
StartLimitIntervalSec=50s
StartLimitBurst=5
@condconfigexists@

[Service]
Environment=DEBUG_LOGGER=--logger=files
EnvironmentFile=-@environment_file@
ExecStart=@sbindir@/sssd -i ${DEBUG_LOGGER}
Type=notify
NotifyAccess=main
Restart=on-abnormal
@capabilities@
SecureBits=noroot noroot-locked
User=@SSSD_USER@
Group=@SSSD_USER@
@supplementary_groups@

[Install]
WantedBy=multi-user.target
