# src/sysv/systemd/sssd-nss.socket.in
[Unit]
Description=SSSD NSS Service responder socket
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
Before=sssd-autofs.socket sssd-pac.socket sssd-pam.socket sssd-ssh.socket sssd-sudo.socket
DefaultDependencies=no
Conflicts=shutdown.target

[Socket]
ExecStartPre=/bin/sh -c "! (pidof -q sssd_nss && pidof sssd_nss | xargs ps -o mntns -p | grep -q `lsns -n -t mnt -o NS -p \\$\\$`)"
ListenStream=@pipepath@/nss
@nss_socket_user_group@

[Install]
WantedBy=sssd.service
