# src/sysv/systemd/sssd-pam.service.in
[Unit]
Description=SSSD PAM Service responder
Documentation=man:sssd.conf(5)
Requires=sssd.service
PartOf=sssd.service
After=sssd.service

[Install]
Also=sssd-pam.socket sssd-pam-priv.socket

[Service]
ExecStart=@libexecdir@/sssd/sssd_pam --uid 0 --gid 0 --debug-to-files

