# src/sysv/systemd/sssd-autofs.socket.in
[Unit]
Description=SSSD AutoFS Service responder socket
Documentation=man:sssd.conf(5)

[Socket]
ListenStream=@pipepath@/autofs

[Install]
WantedBy=sockets.target
