# src/sysv/systemd/sssd.service.in
[Unit]
Description=System Security Services Daemon
# SSSD must be running before we permit user sessions
Before=systemd-user-sessions.service nss-user-lookup.target
Wants=nss-user-lookup.target
StartLimitIntervalSec=50s
StartLimitBurst=5
@condconfigexists@

[Service]
Environment=DEBUG_LOGGER=--logger=files
EnvironmentFile=-@environment_file@
ExecStartPre=+-/bin/chown -f root:root @sssdconfdir@
ExecStartPre=+-/bin/chown -f root:root @sssdconfdir@/sssd.conf
ExecStartPre=+-/bin/chown -f -R root:root @sssdconfdir@/conf.d
ExecStartPre=+-/bin/chown -f -R root:root @sssdconfdir@/pki
ExecStartPre=+-find -type f -iname @dbpath@/*.ldb -exec /bin/chown -f root:root {};
ExecStartPre=+-find -type f -iname @logpath@/*.log -exec /bin/chown -f root:root {};
ExecStart=@sbindir@/sssd -i ${DEBUG_LOGGER}
Type=notify
NotifyAccess=main
Restart=on-abnormal
@capabilities@
SecureBits=noroot noroot-locked
User=root
Group=root
@supplementary_groups@

[Install]
WantedBy=multi-user.target
