# src/sysv/systemd/sssd-pam.socket.in
[Unit]
Description=SSSD PAM Service responder socket
Documentation=man:sssd.conf(5)

[Socket]
ListenStream=@pipepath@/pam

[Install]
WantedBy=sockets.target
