# src/sysv/systemd/sssd-update-confdb.service.in
[Unit]
Description=SSSD Update ConfDB Service

[Install]
Also=sssd-update-confdb.socket

[Service]
ExecStart=@libexecdir@/sssd/sssd_update_confdb --debug-to-files
