# src/sysv/systemd/sssd-nss.service.in
[Unit]
Description=SSSD NSS Service responder
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
RefuseManualStart=true

[Install]
Also=sssd-nss.socket

[Service]
Environment=DEBUG_LOGGER=--logger=files
EnvironmentFile=-@environment_file@
ExecStart=@libexecdir@/sssd/sssd_nss ${DEBUG_LOGGER} --socket-activated
# No capabilities:
CapabilityBoundingSet=
Restart=on-failure
# 'sssd_nss' is special in that it might be used for resolution of 'User='/'Group='/etc,
# and this may cause the service to hang (loop).
# In case SSSD needs to be configured to run as root, avoid adding 'User=root'/'Group=root' explicitly
# here for this reason (use defaults instead).
# In case SSSD was built with support of running under non-root user, there is a special
# handling in 'libnss_sss' and it is allowed to use build time configured user in 'User='/'Group='
@nss_service_user_group@
@supplementary_groups@
