# src/sysv/systemd/sssd-secrets.service.in
[Unit]
Description=SSSD Secrets Service responder
Documentation=man:sssd-secrets(5)
Requires=sssd-secrets.socket
After=sssd-secrets.socket

[Install]
Also=sssd-secrets.socket

[Service]
Environment=DEBUG_LOGGER=--logger=files
ExecStart=@libexecdir@/sssd/sssd_secrets --uid 0 --gid 0 ${DEBUG_LOGGER}
