# src/sysv/systemd/sssd-pam-priv.socket.in
[Unit]
Description=SSSD PAM Service responder private socket
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
BindsTo=sssd-pam.socket
DefaultDependencies=no
Conflicts=shutdown.target

[Socket]
ExecStartPre=@libexecdir@/sssd/sssd_check_socket_activated_responders -r pam
Service=sssd-pam.service
ListenStream=@pipepath@/private/pam
SocketUser=root
SocketGroup=root
SocketMode=0600

[Install]
WantedBy=sssd.service
