# src/sysv/systemd/sssd-pam.service.in
[Unit]
Description=SSSD PAM Service responder
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
RefuseManualStart=true

[Install]
Also=sssd-pam.socket

[Service]
Environment=DEBUG_LOGGER=--logger=files
EnvironmentFile=-@environment_file@
ExecStart=@libexecdir@/sssd/sssd_pam ${DEBUG_LOGGER} --socket-activated
# 'CAP_DAC_READ_SEARCH' is granted as permitted file capability to be elevated to establish GSS API context
CapabilityBoundingSet= CAP_DAC_READ_SEARCH
Restart=on-failure
User=@SSSD_USER@
Group=@SSSD_USER@
@supplementary_groups@
