# src/sysv/systemd/sssd-ssh.socket.in
[Unit]
Description=SSSD SSH Service responder socket
Documentation=man:sssd.conf(5)

[Socket]
ListenStream=@pipepath@/ssh

[Install]
WantedBy=sockets.target
