# src/sysv/systemd/sssd-secrets.socket.in
[Unit]
Description=SSSD Secrets Service responder socket
Documentation=man:sssd-secrets(5)

[Socket]
ListenStream=@runstatedir@/secrets.socket

[Install]
WantedBy=sockets.target
