# src/sysv/systemd/sssd-nss.service.in
[Unit]
Description=SSSD NSS Service responder
Documentation=man:sssd.conf(5)
After=sssd.service
BindsTo=sssd.service
RefuseManualStart=true

[Install]
Also=sssd-nss.socket

[Service]
Environment=DEBUG_LOGGER=--logger=files
EnvironmentFile=-@environment_file@
ExecStart=@libexecdir@/sssd/sssd_nss ${DEBUG_LOGGER} --socket-activated
Restart=on-failure
# Currently SSSD NSS service ('sssd_nss') can't be started under 'sssd' user
# via systemd due to NSS loop when systemd resolves getgrouplist(sssd).
# Hence 'User=' and 'Group=' aren't set (defaults to root).
