# src/sysv/systemd/sssd-update-confdb.socket.in
[Unit]
Description=SSSD Update ConfDB socket

[Socket]
ListenStream=@localstatedir@/run/update-confdb.socket

[Install]
WantedBy=sockets.target
